# shopping-list.service
[Unit]
Description=Shopping List Web App
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/shopping-list
Environment=PATH=/home/pi/shopping-list/venv/bin
EnvironmentFile=/home/pi/shopping-list/.env
ExecStart=/home/pi/shopping-list/venv/bin/python /home/pi/shopping-list/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target